# .config/systemd/user/automount.service
[Unit]
Description=automounter

[Service]
ExecStart=/home/gustav/bin/automounter.sh

[Install]
WantedBy=multi-user.target
